# piplayer.service
[Unit]
Description=piplayer

[Service]
ExecStart=/home/svenstaro/piplayer

[Install]
WantedBy=default.target
